# supervisor.conf
[program:monitor]
command=/usr/bin/python2.7 /home/osmc/Apps/monitor/monitor.py
process_name=%(program_name)s
directory=/home/osmc/Apps/monitor
autostart=true
autorestart=true
redirect_stderr=true